# 0x1A-application_server/gunicorn.service
[Unit]
# Description for better clarity
Description=Gunicorn instance serving the AirBnB Clone application

# Ensure the service starts after the network is available
After=network.target

[Service]
# Run the service under the 'ubuntu' user
User=ubuntu

# Set the working directory for the service
WorkingDirectory=/home/ubuntu/AirBnB_clone_v4

# Define the Gunicorn startup command with 3 workers, binding to 0.0.0.0:5003
ExecStart=gunicorn --workers 3 --bind 0.0.0.0:5003 web_dynamic.2-hbnb:app

# Restart the service automatically if it fails
Restart=always

# Redirect standard output to a log file
StandardOutput=append:/tmp/airbnb-access.log

# Redirect standard error to a log file
StandardError=append:/tmp/airbnb-error.log

# Identifier for syslog messages
SyslogIdentifier=gunicorn

[Install]
# Enable the service to start when the system reaches multi-user.target
WantedBy=multi-user.target
